# site/profile/files/etc/systemd/system/mongod.service
# File managed by Puppet, do not edit manually
.include /usr/lib/systemd/system/mongod.service
[Service]
Type=simple
# open files
LimitNOFILE=128000
# processes/threads
LimitNPROC=128000
